# files/etc/uci-defaults/custom.sh
#!/bin/sh

# 默认ip
uci set network.lan.ipaddr='10.0.0.1'

# 退出状态码
exit 0
